# systemd-services/user/git-reminder.timer
[Unit]
Description=Check for uncommitted git changes daily

[Timer]
OnBootSec=10min
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
