# systemd-services/user/git-reminder.timer
[Unit]
Description=Check for uncommitted git changes every 2 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
