# systemd-services/user/vps-socks.service
[Unit]
Description=SSH SOCKS tunnel to VPS
After=network-online.target
Wants=network-online.target

[Service]
Environment=AUTOSSH_GATETIME=0
Environment=AUTOSSH_POLL=30
UnsetEnvironment=SSH_AUTH_SOCK
ExecStart=/usr/bin/autossh -M 0 -N vps-socks
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
