# systemd-services/system/git-reminder.service
[Unit]
Description=Git Reminder

[Service]
User=maxwell
ExecStart=/usr/bin/git remind status

[Install]
WantedBy=default.target
