# systemd-services/system/git-reminder.timer
[Unit]
Description=Run Git Reminder every hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
