# .config/systemd/user/pass-secret-service.service.d/override.conf
[Service]
# Keep app secrets (Discord tokens, etc.) isolated from the personal
# ~/.password-store tree. Both stores share the same GPG key.
Environment=PASSWORD_STORE_DIR=%h/.password-store-appsecrets
