# .config/systemd/user/pass-secret-service.service.d/override.conf
[Service]
# Keep app secrets (Discord tokens, etc.) isolated from the personal
# ~/.password-store tree. Both stores share the same GPG key.
Environment=PASSWORD_STORE_DIR=%h/.password-store-appsecrets
# pam_gnupg unlocks the key during login, well before this service starts, so the
# preset is normally already cached by then. This delay is cheap insurance against
# the known pam-gnupg / systemd-user-service startup race (cruegge/pam-gnupg README).
ExecStartPre=/usr/bin/sleep 3
